# etc/rc.local
# Run Commands Before shell Startup
# Like Openning Another Terminal And So on
# Start Typing Your Commands here:

# Your Code Here
