# etc/rc.local
source $PODGL/lib/fsh.cfg.sh
source $PODGL/lib/colors

# Run Commands Before shell Startup
# Like Openning Another Terminal And So on
# Start Typing Your Commands here:
source c

# 
